# etc/cron.d/docker-dynamic-hosts
* *     * * *     root   /usr/bin/docker-dynamic-hosts
